# system_files/etc/systemd/system/nix.mount
[Unit]
Description=Nix Package Manager
DefaultDependencies=no
Before=sockets.target
After=ostree-remount.service
BindsTo=var.mount

[Mount]
What=/var/nix
Where=/nix
Options=bind
Type=none

[Install]
WantedBy=local-fs.target
